# bicrypto.service
[Unit]
Description=Bicrypto V5 Application
After=network.target mysql.service mariadb.service redis.service
Wants=mysql.service mariadb.service redis.service

[Service]
Type=simple
# Update these values for your environment
User=runner
Group=runner
# Update this path to your actual installation directory
WorkingDirectory=/home/runner/work/bicrypto/bicrypto
Environment="PATH=/usr/local/bin:/usr/bin:/bin:/home/runner/.local/bin"
Environment="NODE_ENV=production"
ExecStartPre=/bin/sleep 5
# Use the startup script for better reliability
ExecStart=/home/runner/work/bicrypto/bicrypto/start-bicrypto.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bicrypto

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
